# conf/systemd.service
[Unit]
Description=Vantage - self-hosted Fediverse moderation console
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=__APP__
Group=__APP__
WorkingDirectory=__INSTALL_DIR__
# server/static.mjs resolves dist/index.html relative to the install dir, so
# WorkingDirectory must be the install dir.
EnvironmentFile=__INSTALL_DIR__/.env
Environment="PATH=__PATH_WITH_NODEJS__"
ExecStart=__NODEJS_DIR__/node __INSTALL_DIR__/server/index.mjs
Restart=always
RestartSec=10
StandardOutput=append:/var/log/__APP__/__APP__.log
StandardError=append:/var/log/__APP__/__APP__.log
SyslogIdentifier=__APP__

# Hardening. NOTE: no ProtectHome — YunoHost's data_dir lives under
# /home/yunohost.app/__APP__, which ProtectHome would mask from the service
# (the server then can't open its own DB / master key). ProtectSystem=strict
# makes the whole tree read-only, and ReadWritePaths re-opens exactly the two
# dirs the server writes: the data dir (SQLite + master.key) and the log dir.
NoNewPrivileges=true
ProtectSystem=strict
PrivateTmp=true
ReadWritePaths=__DATA_DIR__ /var/log/__APP__

[Install]
WantedBy=multi-user.target
